# deploy/systemd/coturn-video-call.service
[Unit]
Description=coturn STUN/TURN for Video Call
After=network.target

[Service]
Type=simple
User=turnserver
Group=turnserver
AmbientCapabilities=CAP_NET_BIND_SERVICE
ExecStart=/usr/bin/turnserver -c /etc/coturn/video-call.conf
Restart=on-failure
RestartSec=5
MemoryMax=256M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
